# runner/com.sebastian.tgbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sebastian.tgbot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/mondb/Documents/Projects/hunter/runner/run_tgbot.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/mondb</string>
        <key>PATH</key>
        <string>/Users/mondb/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/mondb/Documents/Projects/hunter</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/mondb/Documents/Projects/hunter/runner/telegram_bot.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mondb/Documents/Projects/hunter/runner/telegram_bot.log</string>
</dict>
</plist>
